# artifacts/scheduled_tasks/PROPFIRM_MonitorParseEnrich.xml
<?xml version="1.0" encoding="UTF-16"?>

<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>

    <URI>\PROPFIRM_MonitorParseEnrich</URI>

  </RegistrationInfo>

  <Principals>

    <Principal id="Author">

      <UserId>S-1-5-21-695917052-2651577891-119978033-1001</UserId>

      <LogonType>InteractiveToken</LogonType>

      <RunLevel>HighestAvailable</RunLevel>

    </Principal>

  </Principals>

  <Settings>

    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>

    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>

    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>

    <IdleSettings>

      <Duration>PT10M</Duration>

      <WaitTimeout>PT1H</WaitTimeout>

      <StopOnIdleEnd>true</StopOnIdleEnd>

      <RestartOnIdle>false</RestartOnIdle>

    </IdleSettings>

  </Settings>

  <Triggers>

    <TimeTrigger>

      <StartBoundary>2025-11-10T15:10:44+01:00</StartBoundary>

      <Repetition>

        <Interval>PT15M</Interval>

        <Duration>P365D</Duration>

        <StopAtDurationEnd>true</StopAtDurationEnd>

      </Repetition>

    </TimeTrigger>

  </Triggers>

  <Actions Context="Author">

    <Exec>

      <Command>C:\Users\saint\Documents\PROPFIRM\tools\run_monitor_wrapper.cmd</Command>

    </Exec>

  </Actions>

</Task>
